# conf/supervisor/pusher.ini
[program:pusher]

command=/var/www/pusher.letgo.com/bin/pusher -Dhttp.port=8080 -Dconfig.file=conf/application-staging.conf -Dlogger.resource=logback-prod.xml -Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=16969 -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.authenticate=false
numprocs=1
directory=/var/www/pusher.letgo.com
umask=022
priority=999
autostart=true
autorestart=true
startsecs=3
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=5
user=deployer
redirect_stderr=true

stdout_logfile=/var/log/supervisor/pusher.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=7
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/supervisor/pusher_error.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=7
stderr_capture_maxbytes=1MB

environment=HOME="/var/www/pusher.letgo.com",ENV="dev",JAVA_OPTS='-cp ".:lib/*" -javaagent:lib/org.aspectj.aspectjweaver-1.8.7.jar'
serverurl=AUTO
